# systemd/th-network-monitor-web.service
[Unit]
Description=TH Network Monitor Web
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=thnm
Group=thnm
WorkingDirectory=/opt/th-network-monitor
EnvironmentFile=/etc/th-network-monitor/.env
Environment=THNM_ENV_FILE=/etc/th-network-monitor/.env
ExecStart=/opt/th-network-monitor/.venv/bin/uvicorn app.main:app --host ${APP_HOST} --port ${APP_PORT}
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/var/lib/th-network-monitor /var/log/th-network-monitor
UMask=0077
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
