# systemd/th-network-monitor-worker.service
[Unit]
Description=TH Network Monitor Worker
After=network-online.target th-network-monitor-web.service
Wants=network-online.target

[Service]
Type=simple
User=thnm
Group=thnm
WorkingDirectory=/opt/th-network-monitor
EnvironmentFile=/etc/th-network-monitor/.env
Environment=THNM_ENV_FILE=/etc/th-network-monitor/.env
ExecStart=/opt/th-network-monitor/.venv/bin/python -m monitor.worker
Restart=always
RestartSec=5
NoNewPrivileges=true
AmbientCapabilities=CAP_NET_RAW
CapabilityBoundingSet=CAP_NET_RAW
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/var/lib/th-network-monitor /var/log/th-network-monitor
UMask=0077
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
