# systemd/th-network-monitor-worker.service
[Unit]
Description=TH Network Monitor Worker
After=network.target th-network-monitor-web.service

[Service]
Type=simple
WorkingDirectory=/home/phamsang/Documents/th-network-monitor
EnvironmentFile=/home/phamsang/Documents/th-network-monitor/.env
ExecStart=/home/phamsang/Documents/th-network-monitor/.venv/bin/python -m monitor.worker
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=read-only
ReadWritePaths=/home/phamsang/Documents/th-network-monitor/data /home/phamsang/Documents/th-network-monitor/logs
UMask=0077
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
